# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/tmp/worker.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors


[program:worker1]
command=php worker.php
autorestart=true
startsecs = 0
exitcodes = 0,1,2

[program:worker2]
command=php worker.php
autorestart=true
startsecs = 0
exitcodes = 0,1,2

[program:worker3]
command=php worker.php
autorestart=true
startsecs = 0
exitcodes = 0,1,2

[program:worker4]
command=php worker.php
autorestart=true
startsecs = 0
exitcodes = 0,1,2

[program:worker5]
command=php worker.php
autorestart=true
startsecs = 0
exitcodes = 0,1,2

[program:worker6]
command=php worker.php
autorestart=true
startsecs = 0
exitcodes = 0,1,2

[program:worker7]
command=php worker.php
autorestart=true
startsecs = 0
exitcodes = 0,1,2

[program:worker8]
command=php worker.php
autorestart=true
startsecs = 0
exitcodes = 0,1,2
